# debian/postinst
#!/bin/sh

cp /usr/share/signNativeOsService/turksat.crt /usr/local/share/ca-certificates/

sudo systemctl daemon-reload
sudo systemctl enable turksat-imza.service
sudo systemctl start turksat-imza
